# systemd/vast-host-installer-first-run-notice.service
[Unit]
Description=Vast Host Installer first-run login banner marker
After=network-online.target ssh.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
